# launchd/com.rodrigodapazkaveski.opencode.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.rodrigodapazkaveski.opencode.telegram-bot</string>
    <key>ProgramArguments</key>
    <array>
      <string>/Users/rodrigodapazkaveski/.config/opencode/bin/opencode-telegram-service.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>WorkingDirectory</key>
    <string>/Users/rodrigodapazkaveski</string>
    <key>StandardOutPath</key>
    <string>/Users/rodrigodapazkaveski/Library/Logs/opencode-remote/opencode-telegram.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/rodrigodapazkaveski/Library/Logs/opencode-remote/opencode-telegram.err.log</string>
    <key>ProcessType</key>
    <string>Background</string>
  </dict>
</plist>
